# .devcontainer/supervisord.conf
[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
loglevel=info

[program:jupyter]
# Wait for postCreateCommand to finish creating .venv before starting.
# LD_LIBRARY_PATH is resolved at runtime so it survives Nix store updates.
command=/bin/bash -c "while [ ! -f /workspace/.venv/bin/jupyter-lab ]; do sleep 5; done; export LD_LIBRARY_PATH=\"$(find /nix/store -name libstdc++.so.6 -print -quit 2>/dev/null | xargs -r dirname)\"; exec /workspace/.venv/bin/jupyter-lab --ip=0.0.0.0 --port=8888 --no-browser --allow-root"
directory=/workspace
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/jupyter.log
stderr_logfile=/tmp/jupyter.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=0
